# inference/server/python/cherrypy.conf
[supervisord]

[program:speech-to-text-server]
command=gunicorn wsgi --config /stt/gunicorn.conf --chdir=/stt